# src/main/resources/config.ini
[program:ping1]
command=ping 127.0.0.1
directory=/tmp
stdout_logfile=/tmp/ping.out.log
stderr_logfile=/tmp/ping.err.log
autostart=true
numprocs=3
user=root

[program:ping2]
command=ping 0.0.0.0
directory=/tmp
stdout_logfile=/tmp/ping.out.log
stderr_logfile=/tmp/ping.err.log
autostart=true
numprocs=3
user=root

